# autostart/ssh-add.desktop
[Desktop Entry]
Exec=zsh -extendedglob -c 'ssh-add -q ~/.ssh/ssh^*.pub < /dev/null'
Name=ssh-add
Type=Application
X-KDE-AutoStartScript=true
